# config/profile.d/10-ssh-agent.sh.FreeBSD
# XDG_RUNTIME_DIR guarenteed set by 00-xdg_runtime_dir.sh.FreeBSD
export SSH_AUTH_SOCK="$XDG_RUNTIME_DIR/ssh-agent.socket"
if [ ! -S "$SSH_AUTH_SOCK" ]; then
	eval "$(ssh-agent -a "$SSH_AUTH_SOCK")"
fi
