# config/profile.d/guile-version.sh
# for automake
export GUILE_EFFECTIVE_VERSION=2.2
