# config/profile.d/pycache.sh
export PYTHONPYCACHEPREFIX="${XDG_CACHE_HOME}/python/__pycache__"
