# config/profile.d/setxkbmap.sh.lysator.liu.se
[ -z "$DISPLAY" ] && return
[ -n "$SSH_CLIENT" ] && return

setxkbmap se dvorak_a5 caps:escape
