# config/profile.d/setxkbmap.sh.victoria
[ -z $DISPLAY ] && return

setxkbmap us dvorak
