# config/profile.d/systemd-display-env.sh
[ -z "$DISPLAY" ] && return

systemctl import-environment --user DISPLAY
