# config/profile.d/xinit-xcape.sh.gandalf
[ -z "$DISPLAY" ] && return

xcape -e "Shift_L=parenleft;Shift_R=parenright"
