# config/profile.d/xmonad-state.sh.lysator.liu.se
test -r "$HOME/.xmonad/xmonad.state" && rm "$HOME/.xmonad/xmonad.state"
